# tar/config/.config/autostart/00-env-dbus-propagate.desktop
[Desktop Entry]
Type=Application
Name=Env & DBus propagate
# 로그인 직후, 세션 환경을 DBus activation 경로로 재전파
Exec=/usr/bin/env bash -lc 'command -v dbus-update-activation-environment >/dev/null 2>&1 && dbus-update-activation-environment --all || true'
X-GNOME-Autostart-enabled=true
NoDisplay=true
